# scripts/systemd/script.service.j2

[Unit]
Description="Runs {{ script_path }}"

[Service]
ExecStart="{{ script_path }}"

[Install]
WantedBy=multi-user.target
